# init/gentoo/init.d/libvirt-wakeonlan
#!/sbin/runscript

description="Starts libvirt instances from wake on lan packets"

depend() {
    need net
    after libvirtd
}

start() {
    ebegin "Starting libvirt-wakeonlan"
    start-stop-daemon --start \
	--background -m \
	--pidfile=/var/run/libvirt-wakeonlan.pid \
    -d /usr/share/libvirt-wakeonlan \
    --exec /usr/share/libvirt-wakeonlan/libvirtwol.py -- ${LIBVIRTDWOL_INTERFACE}
    eend $?
}

stop() {

    ebegin "Stopping libvirt-wakeonlan"

    start-stop-daemon --stop --quiet --pidfile=/var/run/libvirt-wakeonlan.pid
    eend $?
}

